# loadbalancer/scripts/rc.local
#!/bin/sh

# start the LoadBalancer
sudo java -jar /home/ec2-user/loadbalancer.jar &
